# build/favicon-fetcher.service
[Unit]
Description=Favicon Fetching API Service
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=0
Restart=yes
ExecStartPre=-/usr/bin/docker stop favicon-fetcher
ExecStartPre=-/usr/bin/docker rm favicon-fetcher
ExecStartPre=/usr/bin/docker build -t favicon-fetcher /usr/src/app
ExecStart=/usr/bin/docker run \
    --name favicon-fetcher \
    -p 8000:8000/udp

[Install]
WantedBy=multi-user.target